# ansible/mediawiki/files/check_master.service
[Unit]
Description=Проверка доступности мастер БД и промоут реплики

[Service]
Type=oneshot
ExecStart=/usr/local/bin/check_master.sh