# script/backup_fs.service
[Unit]
Description=Резервное копирование фс вики

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup_fs.sh
